# dot-config/systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=simple
ExecStart=/home/rgr/bin/getmails

RuntimeMaxSec=480
